# roles/brightnessd/files/brightnessd.service
[Unit]
Description=Joggler Brightness Daemon

[Service]
ExecStart=/usr/sbin/brightnessd
